# BlikCloudFramework/obsolete/init_scripts/cloud-manager-console
#!/sbin/runscript

depend() {
	need syslog-ng
    need postgresql-9.0
}

start() {
	[ -d /var/run/cloud_manager_console ] || mkdir -m 777 /var/run/cloud_manager_console
	ebegin "Starting BlikCloud Console"
	start-stop-daemon --start --quiet --background \
		--exec consoleManager \
		--make-pidfile --pidfile /var/run/cloud_manager_console/cloud_manager_console.pid
	eend $?
}

stop() {
	ebegin "Stopping BlikCloud Console"
	start-stop-daemon --stop --signal INT  \
		--pidfile /var/run/cloud_manager_console/cloud_manager_console.pid
	eend $?
}

reload() {
	if [ ! -f /var/run/cloud_manager_console/cloud_manager_console.pid ] ; then
		eerror "BlikCloud Console is not running."
	else
        stop
	fi

	start
}
